# ddns-gandi.timer
[Unit]
Description=periodically run dns_update.sh

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
